# .clawteam/autostart/com.clawteam.autostart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawteam.autostart</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source ~/.zshrc &amp;&amp; cd /Users/abner/My_Project/HOOK_GZEAMS &amp;&amp; clawteam board live hook-team --interval 5</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/abner/.clawteam/logs/autostart.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/abner/.clawteam/logs/autostart.error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/abner/My_Project/HOOK_GZEAMS</string>

    <key>Environment</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
